# assets/etc/supervisord/apache.conf
[program:apache]
command = httpd -DFOREGROUND
autorestart = false
startentries = 0
stderr_logfile=/proc/1/fd/2
stderr_logfile_maxbytes=0