# assets/etc/supervisord/nginx.conf
[program:nginx]
command = /usr/sbin/nginx
autorestart = false
startentries = 0