# assets/etc/supervisord/nginx.conf
[program:nginx]
command = /usr/sbin/nginx
autorestart = false
startentries = 0
stderr_logfile = /proc/1/fd/2
stderr_logfile_maxbytes = 0