# assets/etc/supervisord.conf
[supervisord]
nodaemon=true

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[unix_http_server]
file=/var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:rtorrent]
command=rtorrent -i %(ENV_WAN_IP)s
user=rtorrent
environment=HOME="/home/rtorrent", PWD="/home/rtorrent", TERM="xterm"
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:rtorrent-logging]
command=tail -f /home/rtorrent/.log/rtorrent.log
user=rtorrent
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
